# infra/launchd/com.daycast.api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.daycast.api</string>

	<key>ProgramArguments</key>
	<array>
		<string>/Users/andrewmaier/daycast/daycast-api/.venv/bin/uvicorn</string>
		<string>app.main:app</string>
		<string>--host</string>
		<string>0.0.0.0</string>
		<string>--port</string>
		<string>8000</string>
		<string>--workers</string>
		<string>2</string>
	</array>

	<key>WorkingDirectory</key>
	<string>/Users/andrewmaier/daycast/daycast-api</string>

	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/Users/andrewmaier/daycast/daycast-api/.venv/bin:/usr/local/bin:/usr/bin:/bin</string>
	</dict>

	<key>RunAtLoad</key>
	<true/>

	<key>KeepAlive</key>
	<true/>

	<key>StandardOutPath</key>
	<string>/Users/andrewmaier/daycast/logs/api-stdout.log</string>

	<key>StandardErrorPath</key>
	<string>/Users/andrewmaier/daycast/logs/api-stderr.log</string>
</dict>
</plist>
